# web-dashboard/aegis-dashboard-optimized.service
[Unit]
Description=Optimized Aegis Security Suite Dashboard
Documentation=https://github.com/aegis-security-suite/dashboard
After=network.target
Wants=network.target

[Service]
Type=simple
User=${CURRENT_USER}
WorkingDirectory=${SECURITY_SUITE_HOME}/web-dashboard
ExecStart=/usr/bin/python3 ${SECURITY_SUITE_HOME}/web-dashboard/app.py
ExecReload=/bin/kill -HUP 
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
Restart=on-failure
RestartSec=30
StartLimitBurst=3
StartLimitInterval=300

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=${SECURITY_SUITE_HOME}/logs
ReadWritePaths=${SECURITY_SUITE_HOME}/configs
ReadWritePaths=${SECURITY_SUITE_HOME}/web-dashboard
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictRealtime=true
MemoryMax=512M
CPUQuota=50%

# Environment
Environment=SECURITY_SUITE_HOME=${SECURITY_SUITE_HOME}
Environment=FLASK_ENV=production
Environment=PYTHONPATH=${SECURITY_SUITE_HOME}/web-dashboard

[Install]
WantedBy=multi-user.target
