# web-dashboard/garuda-dashboard.service
[Unit]
Description=Garuda Security Suite Web Dashboard
Documentation=https://github.com/garuda-security-suite/web-dashboard
After=network.target
Wants=network.target

[Service]
Type=simple
User=${CURRENT_USER}
WorkingDirectory=${SECURITY_SUITE_HOME}/web-dashboard
Environment=SECURITY_SUITE_HOME=${SECURITY_SUITE_HOME}
Environment=FLASK_ENV=production
Environment=FLASK_APP=app-optimized.py
Environment=PYTHONPATH=${SECURITY_SUITE_HOME}/web-dashboard
ExecStart=/usr/bin/python3 ${SECURITY_SUITE_HOME}/web-dashboard/app-optimized.py
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
TimeoutStopSec=30
PrivateTmp=true
Restart=on-failure
RestartSec=10

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=${SECURITY_SUITE_HOME}/configs
ReadWritePaths=${SECURITY_SUITE_HOME}/logs
ReadWritePaths=/tmp
ReadWritePaths=${SECURITY_SUITE_HOME}/web-dashboard

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=garuda-dashboard

[Install]
WantedBy=multi-user.target
Alias=garuda-dashboard.service